# haiku-wifi/init.d/haiku
#!/bin/sh /etc/rc.common

START=99
STOP=15

start() {
        echo start
        cd /root/arthackday
        python haiku.py &
}

stop() {
        echo stop
}
